# nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
stdout_logfile=/var/log/nginx/nginx.log
stderr_logfile=/var/log/nginx/nginx.err

[program:nginx-log-daemon]
command=/usr/local/bin/nginx-log-daemon.sh
autorestart=true
stdout_logfile=/var/log/nginx/log-daemon.log
stderr_logfile=/var/log/nginx/log-daemon.err

[program:certbot-init]
command=/usr/local/bin/init.sh
autorestart=false
stdout_logfile=/var/log/nginx/init.log
stderr_logfile=/var/log/nginx/init.err
